# src/cron/crontab
* * * * * python /opt/batch_import.py
* * * * * python /opt/backup.py