# src/cron/crontab
* * * * * python /opt/batch_import.py
3 * * * * python /opt/backup.py