# src/cron/crontab
* * * * * /opt/file_imports/execute_import.sh